# kpler/src_automatic_09_2025/kpler_fetch.timer
[Unit]
Description=Timer per executar Kpler Fetch cada dia a les 2 AM

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
